# supervisord.conf
[program:celeryd]
command={{ PYTHON }} {{ PROJECT_DIR }}/manage.py celeryd -l info
[program:autoreload]
exclude=true
[program:runserver]
exclude=true
